# blupi-supervisor.conf
[program:blupi]
command = python3 /root/blupi/blupi-scan.py --debug
process_name=%(program_name)s
directory = /root/blupi
user=root
autostart=true
autorestart=true
redirect_stderr=true